# ops/ci-pull-deploy-golang.timer
[Unit]
Description=Poll GitHub for a new Go PM API build every 2 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=2min
AccuracySec=10s

[Install]
WantedBy=timers.target
